# srcs/mysql/supervisor.conf
[supervisord]
nodaemon=true
user=root
[program:mysql]
command = sh -c "mysqld -uroot && kill -s SIGTERM &(cat supervisord.pid)"
autostart=true
autorestart=true
[program:telegraf]
command = sh -c "/etc/telegraf/usr/bin/telegraf --config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"
autostart=true
autorestart=true